# overlay/etc/rc.local
#!/bin/sh -e
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.

su user --login -c startx >/var/log/startx.log 2>&1

exit 0
